# docker/php/crontab
*/5 * * * * /usr/local/bin/php /app/bin/cli.php crawler > /proc/1/fd/1 2>/proc/1/fd/2
* * * * *   sleep 30 && /usr/local/bin/php /app/bin/cli.php notifier > /proc/1/fd/1 2>/proc/1/fd/2
